# templates/jitsi.service.erb
[Unit]
Description=jitsi service with docker compose
Requires=docker.service
After=docker.service

[Service]
User=root
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/srv/jitsi
ExecStart=/usr/local/bin/docker-compose -f docker-compose.yml<% if @compose_jigasi == true -%>
  -f jigasi.yml <% end %><% if @compose_jibri == true -%>
  -f jibri.yml <% end %><% if @compose_etherpad == true -%>
  -f etherpad.yml <% end -%>
  up -d --remove-orphans
ExecStop=/usr/local/bin/docker-compose down

[Install]
WantedBy=multi-user.target

